# pi-factory/systemd/pi-factory-watchdog.service
[Unit]
Description=Pi Factory Watchdog
After=pi-factory.service
Requires=pi-factory.service

[Service]
Type=simple
User=pi
ExecStart=/bin/bash -c 'while true; do sleep 30; curl -sf http://localhost:8000/api/status > /dev/null || (echo "Pi Factory unresponsive, restarting..." | systemd-cat -t pi-factory-watchdog && sudo systemctl restart pi-factory); done'
Restart=always
RestartSec=60

[Install]
WantedBy=multi-user.target
