# pi-factory/systemd/pi-factory.service
[Unit]
Description=Pi Factory — Cosmos Copilot Wizard (v1.1.1)
Documentation=https://github.com/Mikecranesync/factorylm-cosmos-cookoff
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/factorylm-cosmos-cookoff
Environment="PATH=/home/pi/factorylm/services/plc-modbus/.venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="FACTORYLM_NET_MODE=real"
Environment="FACTORYLM_NET_DB=/home/pi/factorylm-cosmos-cookoff/net.db"
ExecStart=/home/pi/factorylm/services/plc-modbus/.venv/bin/uvicorn net.api.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pi-factory

[Install]
WantedBy=multi-user.target
