# pi-factory/configs/pi-factory.service
[Unit]
Description=Pi Factory — Edge Gateway (v3.0)
Documentation=https://github.com/Mikecranesync/factorylm-cosmos-cookoff
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/opt/factorylm
Environment="PATH=/opt/factorylm/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="FACTORYLM_NET_MODE=real"
Environment="FACTORYLM_NET_DB=/opt/factorylm/data/net.db"
Environment="PLC_HOST=192.168.1.100"
Environment="PLC_PORT=502"
Environment="PI_COMPACTCOM_PORT=5020"
# Optional — enable when VFD + camera are wired:
# Environment="VFD_HOST=192.168.1.100"
# Environment="VFD_SLAVE=1"
# Environment="VIDEO_SOURCE=0"
ExecStart=/opt/factorylm/venv/bin/python server.py --port 8000
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pi-factory

[Install]
WantedBy=multi-user.target
